# cgate3/rootfs/etc/s6-overlay/s6-rc.d/cgate/run
#!/command/with-contenv bashio
# shellcheck shell=bash


cd /cgate || bashio::exit.nok "Could not change directory to CGATE"

CGATE_HOME = /config/CGate3/home
CGATE_APPDATA = /config/CGate3/appdata

bashio::log.info "Starting CGATE..."

exec java -jar cgate.jar