# cgate/rootfs/etc/s6-overlay/s6-rc.d/ser2sock/run
#!/command/with-contenv bashio
# shellcheck shell=bash


bashio::log.info "Starting ser2sock..."

exec ser2sock -g 2 -f /etc/ser2sock.conf -p -s $(bashio::config 'serialDevice')
